# SimpleMessageQueue.service
[Unit]
Description=Simple Message Queue Service

[Service]
ExecStart=/usr/bin/python3 /usr/local/lib/python3.8/dist-packages/SimpleMessageQueue.py
